# dist/linux/bin/batch-scheduler.service
[Unit]
Description=batch-scheduler
After=syslog.target network.target zookeeper.service kafka.service

[Service]
Type=forking
Restart=always
RestartSec=1
SuccessExitStatus=143

User=batch
Group=batch

ExecStart=/opt/batch/server/bin/batch-scheduler.sh start
ExecStop=/opt/batch/server/bin/batch-scheduler.sh stop
ExecReload=/opt/batch/server/bin/batch-scheduler.sh reload

[Install]
WantedBy=multi-user.target
